# etc/crontab
# m h  dom mon dow   command
* * * * * ~/bin/ctagger
